# root/etc/services.d/qbittorrent/finish
#!/command/execlineb -S2
define emsg "exit code: ${1} - kill code: ${2}"

ifelse { eltest ${1} -ne 0 }
{
    foreground { echo "*** qBittorrent !! CRASHED !! - ${emsg} - terminating container !!***" }
    /run/s6/basedir/bin/halt
}
foreground { echo "*** qBittorrent finished - ${emsg} ***" }
